# etc/profile.d/makeflags.sh
export CMAKE_BUILD_PARALLEL_LEVEL="`nproc`"
export MAKEFLAGS="-j $CMAKE_BUILD_PARALLEL_LEVEL"

